# etc/init/swapfile.conf
description "swapfile"

start on runlevel [2345]
stop on runlevel [!2345]

script
cd /home/ubuntu/scewpt
exec etc/bin/create_swapfile.sh
end script